# docker-app.service
[Unit]
Description=Docker Application Auto-Start Service
After=docker.service
Requires=docker.service
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/docker-app-start.sh
StandardOutput=journal
StandardError=journal

# Restart policy
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
